# deploy/systemd/quant-lab-risk-permission.service
[Unit]
Description=quant-lab risk permission publisher
After=quant-lab-v5-telemetry-sync.service quant-lab-alpha-evidence.service quant-lab-cost-calibration.service

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
RuntimeMaxSec=5min
MemoryMax=1G
CPUQuota=60%
ExecStart=/opt/quant-lab/.venv/bin/qlab publish-risk-permission --lake-root /var/lib/quant-lab/lake --strategy v5 --version 5.0.0
