# deploy/systemd/quant-lab-alpha-factory.timer
[Unit]
Description=Run quant-lab alpha factory research builder every 6 hours

[Timer]
OnActiveSec=45min
OnUnitActiveSec=6h
RandomizedDelaySec=10min
Persistent=false

[Install]
WantedBy=timers.target
