# deploy/systemd/quant-lab-risk-permission.service
[Unit]
Description=quant-lab risk permission publisher
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=QUANT_LAB_RISK_PERMISSION_TTL_SECONDS=5400
TimeoutStartSec=5min
MemoryMax=1G
CPUQuota=60%
ExecStart=/opt/quant-lab/.venv/bin/qlab publish-risk-permission --lake-root /var/lib/quant-lab/lake --strategy v5 --version 5.0.0
