# deploy/systemd/quant-lab-feature-publish.timer
[Unit]
Description=Run quant-lab feature publisher every 15 minutes

[Timer]
OnActiveSec=45min
OnUnitActiveSec=15min
RandomizedDelaySec=3min
Persistent=false

[Install]
WantedBy=timers.target
