# deploy/systemd/quant-lab-okx-ws.service
[Unit]
Description=quant-lab OKX public WebSocket collector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/quant-lab
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/quant-lab/.venv/bin/qlab okx-ws-collect-universe --symbols BTC-USDT,ETH-USDT,SOL-USDT,BNB-USDT --channels tickers,trades,books5 --market-type SPOT --lake-root /var/lib/quant-lab/lake --flush-interval-seconds 10 --flush-max-messages 100
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
