# deploy/systemd/quant-lab-v5-telemetry-sync.service
[Unit]
Description=quant-lab V5 telemetry bundle sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
ExecStart=/opt/quant-lab/.venv/bin/qlab sync-v5-telemetry --config /etc/quant-lab/v5_telemetry_remote.yaml
RuntimeMaxSec=8min
MemoryHigh=1G
MemoryMax=1536M
NoNewPrivileges=true
PrivateTmp=true
