# deploy/systemd/quant-lab-research-result-import.service
[Unit]
Description=Strictly validate and import signed NAS research results
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
EnvironmentFile=/etc/quant-lab/research-plane.env
Environment=PYTHONUNBUFFERED=1
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
UMask=0007
ExecCondition=/bin/bash -lc '[ "${QUANT_LAB_NAS_RESEARCH_ENABLED:-0}" = "1" ]'
ExecStart=/bin/bash -lc 'set -euo pipefail; COMMIT="$(git -C /opt/quant-lab rev-parse HEAD)"; exec /opt/quant-lab/.venv/bin/qlab import-research-results --lake-root "${QUANT_LAB_LAKE_ROOT}" --queue-root "${QUANT_LAB_RESEARCH_QUEUE_ROOT}" --task-public-key "${QUANT_LAB_RESEARCH_TASK_PUBLIC_KEY_PATH}" --task-key-id "${QUANT_LAB_RESEARCH_TASK_KEY_ID}" --worker-public-key "${QUANT_LAB_RESEARCH_WORKER_PUBLIC_KEY_PATH}" --worker-key-id "${QUANT_LAB_RESEARCH_WORKER_KEY_ID}" --quant-lab-commit "$COMMIT" --max-result-bytes "${QUANT_LAB_RESEARCH_MAX_RESULT_BYTES:-268435456}" --factor-factory-max-result-bytes "${QUANT_LAB_FACTOR_FACTORY_MAX_RESULT_BYTES:-536870912}" --factor-factory-max-value-partition-bytes "${QUANT_LAB_FACTOR_FACTORY_MAX_VALUE_PARTITION_BYTES:-134217728}" --factor-factory-max-value-partition-uncompressed-bytes "${QUANT_LAB_FACTOR_FACTORY_MAX_VALUE_PARTITION_UNCOMPRESSED_BYTES:-268435456}" --factor-factory-max-file-count "${QUANT_LAB_FACTOR_FACTORY_MAX_FILE_COUNT:-20000}" --factor-factory-max-uncompressed-bytes "${QUANT_LAB_FACTOR_FACTORY_MAX_UNCOMPRESSED_BYTES:-1073741824}" --v5-candidate-evidence-max-result-bytes "${QUANT_LAB_V5_CANDIDATE_EVIDENCE_MAX_RESULT_BYTES:-268435456}" --v5-candidate-evidence-max-result-uncompressed-bytes "${QUANT_LAB_V5_CANDIDATE_EVIDENCE_MAX_RESULT_UNCOMPRESSED_BYTES:-536870912}" --v5-candidate-evidence-max-partition-bytes "${QUANT_LAB_V5_CANDIDATE_EVIDENCE_MAX_PARTITION_BYTES:-67108864}" --v5-candidate-evidence-max-partition-uncompressed-bytes "${QUANT_LAB_V5_CANDIDATE_EVIDENCE_MAX_PARTITION_UNCOMPRESSED_BYTES:-134217728}" --v5-candidate-evidence-max-file-count "${QUANT_LAB_V5_CANDIDATE_EVIDENCE_MAX_FILE_COUNT:-5000}" --trade-level-history-max-result-bytes "${QUANT_LAB_TRADE_LEVEL_HISTORY_MAX_RESULT_BYTES:-536870912}" --trade-level-history-max-result-uncompressed-bytes "${QUANT_LAB_TRADE_LEVEL_HISTORY_MAX_RESULT_UNCOMPRESSED_BYTES:-1073741824}" --trade-level-history-max-partition-bytes "${QUANT_LAB_TRADE_LEVEL_HISTORY_MAX_PARTITION_BYTES:-134217728}" --trade-level-history-max-partition-uncompressed-bytes "${QUANT_LAB_TRADE_LEVEL_HISTORY_MAX_PARTITION_UNCOMPRESSED_BYTES:-268435456}" --trade-level-history-max-file-count "${QUANT_LAB_TRADE_LEVEL_HISTORY_MAX_FILE_COUNT:-20000}"'
TimeoutStartSec=30min
Nice=5
CPUQuota=40%
MemoryHigh=2G
MemoryMax=3G
TasksMax=64
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/quant-lab
ReadWritePaths=/var/lib/quant-lab/lake /var/lib/quant-lab/research_queue

[Install]
WantedBy=multi-user.target
