# deploy/systemd/quant-lab-alpha-evidence.timer
[Unit]
Description=Run quant-lab alpha evidence publisher every 15 minutes

[Timer]
OnActiveSec=8min
OnUnitActiveSec=15min
RandomizedDelaySec=1min
Persistent=false

[Install]
WantedBy=timers.target
