# deploy/systemd/quant-lab-feature-publish.service
[Unit]
Description=quant-lab core feature publisher
After=quant-lab-okx-rest-backfill.service

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
TimeoutStartSec=11min
MemoryMax=2G
CPUQuota=80%
ExecStart=/opt/quant-lab/.venv/bin/qlab publish-features --lake-root /var/lib/quant-lab/lake --feature-set core --feature-version v0.1 --timeframe 1H --symbols BTC-USDT,ETH-USDT,SOL-USDT,BNB-USDT
