# deploy/systemd/quant-lab-risk-permission.service
[Unit]
Description=quant-lab risk permission publisher

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
ExecStart=/opt/quant-lab/.venv/bin/qlab publish-risk-permission --lake-root /var/lib/quant-lab/lake --strategy v5 --version 5.0.0
