# deploy/systemd/quant-lab-risk-permission.timer
[Unit]
Description=Run quant-lab risk permission publisher every 10 minutes

[Timer]
OnActiveSec=6min
OnUnitActiveSec=10min
AccuracySec=1min
RandomizedDelaySec=30s
Persistent=false

[Install]
WantedBy=timers.target
