# deploy/systemd/quant-lab-trade-level-history-request.timer
[Unit]
Description=Poll Trade-Level History identity after Candidate Evidence processing

[Timer]
OnCalendar=*-*-* *:50:00 UTC
RandomizedDelaySec=2min
Persistent=true
Unit=quant-lab-trade-level-history-request.service

[Install]
WantedBy=timers.target
