# deploy/systemd/quant-lab-v5-daily-analysis.timer
[Unit]
Description=Run quant-lab V5 telemetry analysis every 5 minutes

[Timer]
OnActiveSec=4min
OnUnitActiveSec=5min
RandomizedDelaySec=45s
Persistent=false

[Install]
WantedBy=timers.target
