# deploy/systemd/quant-lab-lake-compaction.service
[Unit]
Description=Compact quant-lab high-frequency lake datasets
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=PYTHONUNBUFFERED=1
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
Environment=COMPACT_RAW_OKX_WS=1
Environment=COMPACT_HOT_WS_PARTITION_REPAIR=0
Environment=COMPACT_DATASET_TIMEOUT_SECONDS=300
Environment=COMPACT_RUN_BUDGET_SECONDS=1800
Environment=COMPACT_DIRECT_MAX_SOURCE_FILES=64
Environment=COMPACT_DIRECT_MIN_SOURCE_FILES=16
Environment=COMPACT_MAX_SOURCE_BATCH_BYTES=134217728
Environment=COMPACT_CONSOLIDATE_EXISTING_COMPACT_OUTPUTS=0
Environment=MARKET_ROLLUP_LOOKBACK_HOURS=24
Environment=MARKET_ROLLUP_TIMEOUT_SECONDS=600
Environment=MARKET_ROLLUP_POLARS_MAX_THREADS=2
Environment=COMPACT_SMALL_FILE_MAINTENANCE=1
Environment=COMPACT_SMALL_FILE_MAINTENANCE_TIMEOUT_SECONDS=300
Environment=COMPACT_SMALL_FILE_MAINTENANCE_MAX_GROUPS=6
Environment=COMPACT_SMALL_FILE_MAINTENANCE_MAX_SOURCE_FILES_PER_GROUP=64
Environment=COMPACT_SMALL_FILE_MAINTENANCE_TARGET_ROWS=500000
Environment=COMPACT_SMALL_FILE_MAINTENANCE_DATASETS=silver/v5_quant_lab_request,silver/v5_candidate_event,silver/v5_decision_audit,bronze/api_request_metrics
ExecStart=/bin/bash -lc 'set -euo pipefail; /usr/bin/flock -E 75 -w 600 /var/lock/quant-lab-heavy.lock /bin/bash /opt/quant-lab/deploy/scripts/compact_lake_hot_datasets.sh || { code="$?"; if [ "$${code}" = "75" ]; then echo "SKIP_LAKE_COMPACTION_LOCK_BUSY"; exit 0; fi; exit "$${code}"; }'
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
CPUQuota=80%
MemoryHigh=3G
MemoryMax=4G
TimeoutStartSec=35min
NoNewPrivileges=true
PrivateTmp=true
