# deploy/systemd/quant-lab-risk-permission.timer
[Unit]
Description=Run quant-lab risk permission publisher every 5 minutes

[Timer]
OnBootSec=25min
OnUnitActiveSec=5min
RandomizedDelaySec=1min
Persistent=false

[Install]
WantedBy=timers.target
