# deploy/systemd/quant-lab-cost-calibration.timer
[Unit]
Description=Run quant-lab cost calibration every 60 minutes

[Timer]
OnActiveSec=10min
OnUnitActiveSec=60min
RandomizedDelaySec=5min
Persistent=false

[Install]
WantedBy=timers.target
