# deploy/systemd/quant-lab-alpha-evidence.timer
[Unit]
Description=Run quant-lab alpha evidence publisher hourly

[Timer]
OnBootSec=20min
OnUnitActiveSec=1h
RandomizedDelaySec=3min
Persistent=false

[Install]
WantedBy=timers.target
