# deploy/systemd/quant-lab-okx-rest-backfill.timer
[Unit]
Description=Run quant-lab OKX public REST backfill every hour

[Timer]
OnActiveSec=5min
OnUnitActiveSec=1h
RandomizedDelaySec=2min
Persistent=false

[Install]
WantedBy=timers.target
