# deploy/systemd/quant-lab-alpha-evidence.service
[Unit]
Description=quant-lab alpha evidence publisher
After=quant-lab-feature-publish.service quant-lab-cost-calibration.service

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
RuntimeMaxSec=15min
MemoryHigh=1G
MemoryMax=2G
CPUQuota=80%
ExecStart=/opt/quant-lab/.venv/bin/qlab build-alpha-evidence --lake-root /var/lib/quant-lab/lake --alpha-id v5.core.momentum --version v0.1 --feature-set core --feature-version v0.1 --feature-names close_return_24 --timeframe 1H --label-horizon-bars 4 --decision-delay-bars 1 --universe-id okx-major-spot --strategy v5 --cost-quantile p75
