# deploy/systemd/quant-lab-v5-telemetry-sync.timer
[Unit]
Description=Run quant-lab V5 telemetry sync every 3 minutes

[Timer]
OnActiveSec=2min
OnUnitActiveSec=3min
RandomizedDelaySec=30s
Persistent=false

[Install]
WantedBy=timers.target
