# deploy/systemd/quant-lab-web-v2-smoke.timer
[Unit]
Description=Run quant-lab Web V2/API smoke every 10 minutes

[Timer]
OnActiveSec=3min
OnUnitActiveSec=10min
AccuracySec=1min
RandomizedDelaySec=60s
Persistent=false

[Install]
WantedBy=timers.target
