# deploy/systemd/quant-lab-factor-factory-request.service
[Unit]
Description=Seal full-history Factor Factory Snapshot and enqueue NAS research task
After=quant-lab-feature-publish.service quant-lab-v5-telemetry-sync.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
EnvironmentFile=/etc/quant-lab/research-plane.env
Environment=PYTHONUNBUFFERED=1
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
UMask=0007
ExecCondition=/bin/bash -lc '[ "${QUANT_LAB_NAS_RESEARCH_ENABLED:-0}" = "1" ] && [ "${QUANT_LAB_NAS_FACTOR_FACTORY_ENABLED:-0}" = "1" ]'
ExecStartPre=/usr/bin/test -d /var/lib/quant-lab/lake/bronze/lake_file_index
ExecStartPre=/usr/bin/test -d /var/lib/quant-lab/lake/gold/feature_value
ExecStartPre=/usr/bin/test -d /var/lib/quant-lab/lake/silver/market_bar
ExecStart=/bin/bash -lc 'set -euo pipefail; COMMIT="$(git -C /opt/quant-lab rev-parse HEAD)"; exec /opt/quant-lab/.venv/bin/qlab request-factor-factory --lake-root "${QUANT_LAB_LAKE_ROOT}" --queue-root "${QUANT_LAB_RESEARCH_QUEUE_ROOT}" --key-id "${QUANT_LAB_RESEARCH_TASK_KEY_ID}" --quant-lab-commit "$COMMIT" --date auto --feature-set core --feature-version v0.1 --factor-version v0.1 --timeframe 1H --horizon-bars 4,8,24,72 --decision-delay-bars 1 --max-factors 200 --min-samples 100 --top-quantile 0.2 --cost-quantile p75 --max-input-bytes "${QUANT_LAB_FACTOR_FACTORY_MAX_SNAPSHOT_BYTES:-26843545600}" --max-pending-tasks "${QUANT_LAB_FACTOR_FACTORY_MAX_PENDING_TASKS:-1}"'
TimeoutStartSec=20min
Nice=10
IOSchedulingClass=idle
CPUQuota=30%
MemoryHigh=500M
MemoryMax=900M
TasksMax=64
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadOnlyPaths=/opt/quant-lab /var/lib/quant-lab/lake
ReadWritePaths=/var/lib/quant-lab/research_queue /var/lib/quant-lab/lake/bronze/lake_file_index

[Install]
WantedBy=multi-user.target
