# deploy/systemd/quant-lab-storage-retention.timer
[Unit]
Description=Run quant-lab storage retention pruning daily

[Timer]
OnCalendar=*-*-* 09:20:00
RandomizedDelaySec=10min
Persistent=false

[Install]
WantedBy=timers.target
