# deploy/systemd/quant-lab-okx-rest-backfill.service
[Unit]
Description=quant-lab OKX public REST market-bar backfill
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'set -euo pipefail; for symbol in BTC-USDT ETH-USDT SOL-USDT BNB-USDT; do /opt/quant-lab/.venv/bin/qlab okx-fetch-candles --inst-id "$symbol" --bar 1H --market-type SPOT --lake-root /var/lib/quant-lab/lake --history --limit 100; done'
