# deploy/systemd/quant-lab-daily-export.timer
[Unit]
Description=Run quant-lab daily expert export

[Timer]
OnCalendar=*-*-* 00:20:00
Persistent=true

[Install]
WantedBy=timers.target
