# deploy/systemd/quant-lab-risk-permission.timer
[Unit]
Description=Run quant-lab risk permission publisher every 10 minutes

[Timer]
OnBootSec=3min
OnActiveSec=10min
OnUnitActiveSec=10min
RandomizedDelaySec=30s
Persistent=true

[Install]
WantedBy=timers.target
