# deploy/systemd/quant-lab-okx-ws.service
[Unit]
Description=quant-lab OKX public WebSocket collector
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=PYTHONUNBUFFERED=1
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
Environment=QUANT_LAB_WS_APPEND_TARGET_ROWS=500000
Environment=QUANT_LAB_WS_APPEND_PARTITIONED=0
ExecStart=/opt/quant-lab/.venv/bin/qlab okx-ws-collect-universe --symbols BTC-USDT,ETH-USDT,SOL-USDT,BNB-USDT --channels tickers,trades,books5 --market-type SPOT --lake-root /var/lib/quant-lab/lake --flush-interval-seconds 300 --flush-max-messages 50000
Restart=on-failure
RestartSec=30
MemoryHigh=2G
MemoryMax=3G
KillSignal=SIGINT
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
