# deploy/systemd/quant-lab-daily-export.service
[Unit]
Description=quant-lab daily expert snapshot export
After=quant-lab-risk-permission.service quant-lab-v5-telemetry-sync.service

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
PermissionsStartOnly=true
TimeoutStartSec=45min
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
Environment=QUANT_LAB_EXPORT_V5_MAX_PENDING_BUNDLES=1
Environment=QUANT_LAB_EXPORT_V5_MAX_SCAN_BUNDLES=30
EnvironmentFile=-/etc/quant-lab/quant_lab_api.env
MemoryHigh=2G
MemoryMax=3G
CPUQuota=80%
ExecStartPre=/usr/bin/systemctl start quant-lab-v5-telemetry-sync.service
ExecStartPre=/opt/quant-lab/.venv/bin/qlab publish-risk-permission --lake-root /var/lib/quant-lab/lake --strategy v5 --version 5.0.0
ExecStart=/bin/bash -lc 'set -euo pipefail; /usr/bin/flock -E 75 -w 600 /var/lock/quant-lab-heavy.lock /usr/bin/flock -E 75 -w 600 /var/lock/quant-lab-v5-telemetry-sync.lock /usr/bin/timeout 20m /opt/quant-lab/.venv/bin/qlab export-daily --date "$(date +%%F)" --lake-root /var/lib/quant-lab/lake --out-dir /var/lib/quant-lab/exports --no-refresh-risk-permission --pre-export-v5-refresh --allow-stale-v5 || { code="$$?"; if [ "$${code}" = "75" ]; then echo "SKIP_DAILY_EXPORT_LOCK_BUSY"; exit 0; fi; exit "$${code}"; }'
