# deploy/systemd/quant-lab-lake-compaction.service
[Unit]
Description=Compact quant-lab high-frequency lake datasets
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=PYTHONUNBUFFERED=1
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
ExecStart=/usr/bin/flock -w 10 /var/lock/quant-lab-heavy.lock /bin/bash /opt/quant-lab/deploy/scripts/compact_lake_hot_datasets.sh
Nice=15
IOSchedulingClass=idle
MemoryHigh=5G
MemoryMax=6G
TimeoutStartSec=35min
NoNewPrivileges=true
PrivateTmp=true
