# deploy/systemd/quant-lab-v5-daily-analysis.service
[Unit]
Description=quant-lab V5 telemetry analysis
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
RuntimeMaxSec=10min
MemoryMax=2G
CPUQuota=80%
ExecStart=/opt/quant-lab/.venv/bin/qlab analyze-v5-telemetry --lake-root /var/lib/quant-lab/lake
NoNewPrivileges=true
PrivateTmp=true
