# deploy/systemd/quant-lab-alpha-factory.timer
[Unit]
Description=Run quant-lab alpha factory research builder daily off peak

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=20min
Persistent=false

[Install]
WantedBy=timers.target
