# deploy/systemd/quant-lab-v5-telemetry-sync.service
[Unit]
Description=quant-lab V5 telemetry bundle sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
Environment=QUANT_LAB_V5_SYNC_MAX_BUNDLES=1
Environment=QUANT_LAB_V5_SYNC_REMOTE_MAX_FILES=1
Environment=QUANT_LAB_V5_SYNC_MAX_SCAN_BUNDLES=1
ExecStart=/bin/bash -lc 'set -euo pipefail; /usr/bin/flock -E 75 -w 10 /var/lock/quant-lab-heavy.lock /usr/bin/flock -E 75 -w 10 /var/lock/quant-lab-v5-telemetry-sync.lock /usr/bin/timeout 20m /opt/quant-lab/.venv/bin/qlab sync-v5-telemetry --config /etc/quant-lab/v5_telemetry_remote.yaml --max-bundles 1 --remote-max-files 1 --max-scan-bundles 1 --newest-first --max-skipped-files-reported 20 --skip-historical-outcomes --skip-analysis-after-sync --compact-output || { code="$?"; if [ "$${code}" = "75" ]; then echo "SKIP_V5_TELEMETRY_SYNC_LOCK_BUSY"; exit 0; fi; exit "$${code}"; }'
TimeoutStartSec=22min
MemoryHigh=5G
MemoryMax=6G
CPUQuota=80%
NoNewPrivileges=true
PrivateTmp=true
