# deploy/systemd/quant-lab-v5-telemetry-sync.timer
[Unit]
Description=Run quant-lab V5 telemetry sync every 10 minutes

[Timer]
OnActiveSec=2min
OnUnitActiveSec=10min
RandomizedDelaySec=1min
Persistent=false

[Install]
WantedBy=timers.target
