# deploy/systemd/quant-lab-v5-research-refresh.timer
[Unit]
Description=Run quant-lab V5 candidate research refresh every hour

[Timer]
OnCalendar=*-*-* 00:05:00
OnActiveSec=10min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=false

[Install]
WantedBy=timers.target
