# deploy/systemd/quant-lab-lake-compaction.timer
[Unit]
Description=Run quant-lab lake compaction every hour

[Timer]
OnActiveSec=45min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=false

[Install]
WantedBy=timers.target
