# deploy/systemd/quant-lab-v5-daily-analysis.timer
[Unit]
Description=Run quant-lab V5 telemetry analysis every 30 minutes

[Timer]
OnActiveSec=30min
OnUnitActiveSec=30min
RandomizedDelaySec=5min
Persistent=false

[Install]
WantedBy=timers.target
