# deploy/systemd/quant-lab-v5-daily-analysis.service
[Unit]
Description=quant-lab V5 telemetry analysis
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
RuntimeMaxSec=8min
MemoryHigh=1G
MemoryMax=2G
CPUQuota=60%
ExecStart=/usr/bin/flock -w 30 /var/lock/quant-lab-v5-research.lock /opt/quant-lab/.venv/bin/qlab analyze-v5-telemetry --lake-root /var/lib/quant-lab/lake --skip-candidate-gold
NoNewPrivileges=true
PrivateTmp=true
