# deploy/systemd/quant-lab-alpha-factory.timer
[Unit]
Description=Legacy manual-only local Alpha Factory timer (disabled by default)

[Timer]
OnCalendar=*-*-* 03:20:00
RandomizedDelaySec=20min
Persistent=false

[Install]
WantedBy=timers.target
