# deploy/systemd/quant-lab-okx-rest-backfill.timer
[Unit]
Description=Run quant-lab OKX public REST backfill every 15 minutes

[Timer]
OnActiveSec=15min
OnUnitActiveSec=15min
RandomizedDelaySec=2min
Persistent=false

[Install]
WantedBy=timers.target
