# deploy/systemd/quant-lab-export-request.path
[Unit]
Description=Watch for lightweight NAS Expert Export requests

[Path]
DirectoryNotEmpty=/var/lib/quant-lab/export_queue/requests/pending
Unit=quant-lab-export-request.service

[Install]
WantedBy=multi-user.target
