# deploy/systemd/quant-lab-risk-permission.timer
[Unit]
Description=Run quant-lab risk permission publisher every 3 minutes

[Timer]
OnActiveSec=10min
OnUnitActiveSec=3min
RandomizedDelaySec=30s
Persistent=false

[Install]
WantedBy=timers.target
