# deploy/systemd/quant-lab-okx-readonly-backfill.timer
[Unit]
Description=Run quant-lab OKX read-only private backfill every 15 minutes

[Timer]
OnBootSec=7min
OnUnitActiveSec=1h
RandomizedDelaySec=5min
Persistent=false

[Install]
WantedBy=timers.target
