# deploy/systemd/quant-lab-v5-regime-router.service
[Unit]
Description=quant-lab V5 regime router refresh
After=network-online.target quant-lab-v5-daily-analysis.service quant-lab-v5-research-refresh.service
Wants=network-online.target

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
Environment=POLARS_MAX_THREADS=1
Environment=MALLOC_ARENA_MAX=2
TimeoutStartSec=12min
MemoryHigh=768M
MemoryMax=1G
CPUQuota=40%
ExecStart=/bin/bash -lc '/usr/bin/flock -E 75 -w 5 /var/lock/quant-lab-v5-regime-router.lock /usr/bin/flock -E 75 -w 120 /var/lock/quant-lab-v5-research.lock /usr/bin/timeout 10m /opt/quant-lab/.venv/bin/qlab build-regime-router --lake-root /var/lib/quant-lab/lake --date auto || { code="$?"; if [ "$${code}" = "75" ]; then echo "SKIP_V5_REGIME_ROUTER_LOCK_BUSY"; exit 0; fi; exit "$${code}"; }'
NoNewPrivileges=true
PrivateTmp=true
