# deploy/systemd/quant-lab-cost-calibration.service
[Unit]
Description=quant-lab daily cost calibration
After=quant-lab-okx-rest-backfill.service quant-lab-okx-readonly-backfill.service

[Service]
Type=oneshot
User=quantlab
Group=quantlab
WorkingDirectory=/opt/quant-lab
RuntimeMaxSec=10min
MemoryMax=1G
CPUQuota=80%
ExecStart=/bin/bash -lc '/opt/quant-lab/.venv/bin/qlab calibrate-costs --lake-root /var/lib/quant-lab/lake --day "$(date -u +%%F)" --min-sample-count 30'
