# deploy/systemd/quant-lab-alpha-factory-request.timer
[Unit]
Description=Request NAS Alpha Factory research daily off peak

[Timer]
OnCalendar=*-*-* 03:20:00 UTC
RandomizedDelaySec=20min
Persistent=true
Unit=quant-lab-alpha-factory-request.service

[Install]
WantedBy=timers.target
